# supervisord.conf
[supervisord]
nodaemon=true

[program:bitlbee]
command=/usr/local/sbin/bitlbee -F
